# systemd/fox-bbs.service
[Unit]
Description=Fox BBS - Amateur Radio Bulletin Board System
Documentation=https://github.com/morria/fox
After=network.target network-online.target
Wants=network-online.target

[Service]
Type=simple
User=USER_PLACEHOLDER
Group=GROUP_PLACEHOLDER
WorkingDirectory=INSTALL_DIR_PLACEHOLDER

# Environment
Environment="PATH=VENV_PATH_PLACEHOLDER/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"

# Command to run
ExecStart=VENV_PATH_PLACEHOLDER/bin/python3 INSTALL_DIR_PLACEHOLDER/fox_bbs.py --config CONFIG_PATH_PLACEHOLDER

# Restart policy - restart on failure, but not too aggressively
Restart=on-failure
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=INSTALL_DIR_PLACEHOLDER

# Logging - all output goes to systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fox-bbs

[Install]
WantedBy=multi-user.target
